# Torrent/container-startup.service
[Unit]
Description=Start Docker containers with TUN device
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/root/Git/Containers/Torrent/startup-automation.sh
RemainAfterExit=yes
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 